# data/opt-rootfs-sys.mount
[Unit]
After=opt-rootfs.mount

[Mount]
What=sys
Where=/opt/rootfs/sys
Type=proc
